# claude-mon.service
# systemd *user* service for claude-mon — no root required.
# Assumes you cloned to ~/claude-mon; edit the two paths below if not.
#
#   mkdir -p ~/.config/systemd/user
#   cp claude-mon.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now claude-mon
#   loginctl enable-linger "$USER"     # keep running after logout / across reboots

[Unit]
Description=claude-mon — live Claude Code usage dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/claude-mon
ExecStart=/usr/bin/env python3 %h/claude-mon/server.py --port 8420
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
